# pkgs/50-pipewire/pipewire.install
post_install() {
  # Enable socket by default
  systemctl --global enable pipewire.socket
}

post_upgrade() {
  if (( $(vercmp $2 0.1.6-2) < 0)); then
    systemctl --global enable pipewire.socket
  fi
}

pre_remove() {
  systemctl --global disable pipewire.socket
}
