# rootfs/etc/systemd/system/cec-active-source.timer
[Unit]
Description=Trigger cec-active-source at boot

[Timer]
OnBootSec=1
OnStartupSec=1

[Install]
WantedBy=timers.target
